# deploy/grav-charge-sim.service
[Unit]
Description=Gravitational Charge Simulation (gsim.vdisknow.com)
After=network.target

[Service]
Type=simple
Environment=DOMAIN=gsim.vdisknow.com
Environment=CERT_DIR=/opt/grav-charge-sim/certs
ExecStart=/opt/grav-charge-sim/grav-charge-sim
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
